# services/gcloud-auth.service
[Unit]
Description=google cloud authentication service

[Service]
ExecStart=/usr/local/bin/gcloud-auth.sh

[Install]
WantedBy=default.target
